# playbook/roles/deploying/templates/supervisor.conf
[program:deploying]
command = /home/deploying/deploying-linux-amd64
directory = /home/deploying
autorestart = true
user = deploying
